# tools/rootfs/var/lib/opkg/info/arago-gpl-notice.postrm
#!/bin/sh
update-rc.d $D gplv3-notice remove
