# tools/rootfs/var/lib/opkg/info/initscript-telnetd.postrm
#!/bin/sh
update-rc.d $D telnetd remove
