# tools/rootfs/var/lib/opkg/info/psplash.list
/mnt/.psplash/
/usr/bin/psplash-write
/etc/init.d/psplash.sh
