# tools/rootfs/var/lib/opkg/info/psplash.postrm
#!/bin/sh
update-rc.d $D psplash.sh remove
